# tools/website/crontab
0 4 * * * root /root/code/TopicLand/tools/website/reset_website.sh